# overlay/etc/systemd/system/getty@.service.d/override.conf
[Unit]
After=expand-rootfs.service gen-dropbear-keys.service
