# ship/etc/systemd/system/swapfile.service
[Unit]
Description=Generate swapfile

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/swapfile

[Install]
WantedBy=multi-user.target
