# overlay/etc/systemd/system/wlan-gpio.service
[Unit]
Description=WLAN GPIO

[Service]
Type=simple
ExecStart=/usr/bin/gpioset -c gpiochip0 2=1

[Install]
WantedBy=multi-user.target
