# ship/etc/systemd/system/systemd-firstboot.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/systemd-firstboot --prompt-timezone --prompt-hostname --prompt-root-password --force
